# root/etc/s6-overlay/s6-rc.d/init-doplarr-config/up
/etc/s6-overlay/s6-rc.d/init-doplarr-config/run
